# scripts/openrc/droid4-battery-calibration
#!/sbin/openrc-run

name="droid4-battery-calibrate"

depend() {
	after kmod
}

start() {
	cat /var/lib/droid4-battery-calibration/charge_full > /sys/class/power_supply/battery/charge_full || true
}

stop() {
	ebegin "Saving charge_full"
	mkdir -p /var/lib/droid4-battery-calibration
	cat /sys/class/power_supply/battery/charge_full > /var/lib/droid4-battery-calibration/charge_full || true
}
