# scripts/openrc/maserati-battery-calibration
#!/sbin/openrc-run

name="maserati-battery-calibrate"

depend() {
       after kmod
}

start() {

    if [ -f /var/lib/maserati-battery-calibration/charge_full ]; then
	ebegin "Restoring charge_full"
	cat /var/lib/maserati-battery-calibration/charge_full > /sys/class/power_supply/battery/charge_full
    else
        einfo "No charge_full file to restore"
    fi
}

stop() {
    ebegin "Saveing charge_full"
    if [ ! -d /var/lib/maserati-battery-calibration ]; then
        mkdir /var/lib/maserati-battery-calibration
    fi
    cat /sys/class/power_supply/battery/charge_full > /var/lib/maserati-battery-calibration/charge_full
    if [ $? -ne 0 ]; then
        einfo "Saveing charge_full failed"
        rm /var/lib/maserati-battery-calibration/charge_full
    fi
}


